# .crontab
0 0 * * * /usr/local/bin/remind /Users/dnorris/drive/dcn/remind/defer.rem | /usr/bin/egrep -v '^Reminders for|^$' >> /Users/dnorris/drive/dcn/in.otl
